# autostart/autostart.service
[Unit]
Description=AutoExec
[Service]
ExecStart=/etc/autostart.sh
[Install]
WantedBy=multi-user.target
